# crontab
 0 0 * * * root python /data_remover.py > /proc/1/fd/1 2>/proc/1/fd/2
